# recipes-domd/domd-image-weston/files/meta-xt-prod-extra/recipes-extended/guest-addons/files/android-disks.service
[Unit]
Description=Android domain disks setup
Requires=proc-xen.mount
After=proc-xen.mount
After=network-online.target
After=remote-fs.target

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/xt/scripts/android-disks.sh

[Install]
WantedBy=multi-user.target
